# index.service
[Unit]
Description=My Application Service
After=network.target

[Service]
ExecStart=/opt/index.js

[Install]
WantedBy=multi-user.target